# crontab
@reboot sh -c 'cd ~/meta-lang.com && pm2 start server.mjs'